# tonyPi/FYP_Robot/fyp-robot.service
[Unit]
Description=FYP Robot Main Controller (Vision + Voice + MQTT)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/FYP_Robot
ExecStart=/usr/bin/python3 /home/pi/FYP_Robot/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=MQTT_BROKER=192.168.149.100
Environment=MQTT_PORT=1883
Environment=ROBOT_ID=tonypi_fyp
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
